# systemd/disable-gpu.service
[Unit]
Description=Disable the dedicated GPU

[Service]
Type=oneshot
ExecStart=/usr/local/bin/disable-gpu.sh

[Install]
WantedBy=multi-user.target
